# ops/systemd/patchwork.service
[Unit]
Description=Patchwork Service
After=network.target

[Service]
Type=simple
User=patchwork
Group=patchwork
WorkingDirectory=/opt/patchwork
Environment=PATCHWORK_BIND_ADDR=127.0.0.1:8080
EnvironmentFile=/etc/patchwork/patchwork.env
ExecStart=/opt/patchwork/patchwork
Restart=on-failure
RestartSec=5s
LimitNOFILE=65535
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectControlGroups=true
ProtectKernelModules=true
ProtectKernelTunables=true
RestrictSUIDSGID=true
LockPersonality=true

[Install]
WantedBy=multi-user.target
